# adgen-backend.service
[Unit]
Description=Ad Generation Backend Service
After=network.target

[Service]
User=azureuser
WorkingDirectory=/home/azureuser/marketing-app-ad-gen/backend
Environment="PATH=/home/azureuser/.conda/envs/segp-env/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="CREW_API_KEY=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"
Environment="OPENAI_API_KEY=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"
ExecStart=/home/azureuser/.conda/envs/segp-env/bin/python -m uvicorn main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target 